# .config/local/supervisor.conf
[program:nginx]
command : nginx

[program:uwsgi]
command = uwsgi -i /srv/airbnb/.config/local/uwsgi.ini